# client/ubuntu-desktop-env/conf.d/supervisord.conf
[program:clipboard_app]
priority=25
directory=/root/clipboard_app
command=/usr/local/bin/python /root/clipboard_app/clipboard_server.py
